# deployment/systemd/bot-update-staging.service
[Unit]
Description=Discord Bot Staging Auto-Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root

# bot-update-staging スクリプトを実行
ExecStart=/usr/local/bin/bot-update-staging

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-update-staging

[Install]
WantedBy=multi-user.target
